# lessons/Linux/05-processes-systemd-services/scripts/units/hello.service
[Unit]
Description=Hello logger (oneshot)

[Service]
Type=oneshot
ExecStart=/usr/local/bin/hello.sh
